# crontab
# lame hack to assure that things keep running
*/5 * * * * cd ~ ; make start ; exit 0
